# data/config/s6/services.d/samba/run
#!/usr/bin/with-contenv /bin/bash
# vim:ft=sh

smbd="$(which smbd)"
exec env -i "$smbd" --no-process-group -FD --log-stdout
